# home-assistant-addon/audioserver/rootfs/etc/services.d/lox-audioserver/run
#!/command/with-contenv bashio

set -euo pipefail

cd /opt/lox-audioserver
bashio::log.info "Starting Loxone Audio Server"
exec node dist/server.js
